# conf/supervisor.ini
[program:simbock-queue]
process_name = %(program_name)s
command = /opt/php74/bin/php /root/simbock/artisan queue:listen
autorestart = true
user = root
redirect_stderr = true
stdout_logfile = /root/simbock/storage/logs/workers/queue.log
stopwaitsecs = 3600

[program:simbock-websockets]
process_name = %(program_name)s
command = /opt/php74/bin/php /root/simbock/artisan websockets:serve --host=195.54.162.63
autorestart = true
user = root
redirect_stderr = true
stdout_logfile = /root/simbock/storage/logs/workers/websockets.log
stopwaitsecs = 3600
